# aws deployment scripts/raft.service
[Unit]
Description=Power-off gpu

[Service]
Type=oneshot
ExecStart=$HOME/raft_node_start

[Install]
WantedBy=multi-user.target